# bin/db_puma.service
[Unit]
Description=destinationbook local server management
After=network.target


[Service]
ExecStart=/var/www/destinationbook.com/bin/db_puma.sh start
ExecStop=/var/www/destinationbook.com/bin/db_puma.sh stop
Restart=always
RestartSec=10                       # Restart service after 10 seconds if node service crashes
StandardOutput=syslog               # Output to syslog
StandardError=syslog                # Output to syslog
SyslogIdentifier=jpm

[Install]
WantedBy=multi-user.target
